# systemd/donna.service
[Unit]
Description=Donna stack ensure (snapshot + compose up)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
# Set this to the Discord debug webhook before enabling, or boot alerts are silent.
Environment=DONNA_ALERT_WEBHOOK=
ExecStart=/mnt/donna/donna/scripts/donna-deploy.sh ensure
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
